# kali-anonsurf-deb-src/etc/init/pandora.conf
description "Wipes RAM memory on shutdown"

start on (runlevel [016] and stopped lightdm)

script
	/usr/bin/sdmem -flv
end script
